# crontab
* * * * * zsh ‾/bin/bckup
0 * * * * rm ‾/screenshots/*
0,30 * * * * zsh ‾/bin/sync-dotfiles
